# startagent.service
[Unit]
Description=My script

[Service]
ExecStart=/home/root/iotkit-agent/start-agent.sh

[Install]
WantedBy=multi-user.target 

